# @darwin/home/Library/LaunchAgents/co.joshwel.lottedisplayfriend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">

<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>co.joshwel.lottedisplayfriend</string>

    <key>ProgramArguments</key>
    <array>
      <string>/run/current-system/sw/bin/python</string>
      <string>/Users/majo/Works/Interchange/LotteDisplayFriend/python/lottedisplayfriend.py</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/lottedisplayfriend.out.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/lottedisplayfriend.err.log</string>
  </dict>
</plist>
